# build/supervisord.conf
[supervisord]
nodaemon=true

[program:redis]
command=/usr/bin/redis-server

[program:notorious]
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
command=/go/bin/notorious

[program:mysql]
command=/usr/bin/mysqld_safe
